# server/config/latia.com.ubuntu.conf
#sudo nano /etc/init/latia.com.conf
description "latia.com"
env APP_PATH="/var/www/latia"
start on startup
stop on shutdown
script
  cd $APP_PATH
  exec node app.js
end script